# roles/duplicity/templates/duplicity-backup.timer.j2
# {{ ansible_managed }}

[Unit]
Description=duplicity backup timer
Wants=network-online.target

[Timer]
OnCalendar=*-*-* *:00:00
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
